# pimon_old/raspberry-pi-dashboard/config/autostart/startchrome.desktop
[Desktop Entry] 

Type=Application

Exec=sh /home/pi/.config/autostart/readychrome.sh
